# provisioning/roles/base/files/athinfod.socket
[Unit]
Description=athinfod

[Socket]
ListenStream=49155
Accept=yes        

[Install]
WantedBy=sockets.target
